# 99-switch-display.rules
KERNEL=="card0", SUBSYSTEM=="drm", ENV{DISPLAY}=":0", ENV{XAUTHORITY}="/home/guts/.Xauthority", RUN+="/home/guts/scripts/monitors_fix.sh"
